# data/etc/supervisord.conf
; http://supervisord.org/configuration.html#file-format
; Open server on a UNIX socket
[unix_http_server]
file=/tmp/supervisord.sock
chmod=0700
chown=nobody:nogroup

[supervisord]
; Set up our logs
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
; Keep pidfile same as `/etc/init.d`
pidfile=/var/run/supervisord.pid

[supervisorctl]
; Use our UNIX HTTP server as our server
serverurl=unix:///tmp/supervisord.sock

; https://github.com/Supervisor/supervisor/blob/3.2.0/supervisor/skel/sample.conf#L41-L45
[rpcinterface:supervisor]
; Define protocol for communicating with supervisor server
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; Configure our services
; DEV: We use `-server` suffix to allow for multiple variations per service repo
;   For example, `-celery` and `-redis`
[program:twolfson.com-server]
command=/usr/bin/node /home/ubuntu/twolfson.com/main/bin/twolfson.com
directory=/home/ubuntu/twolfson.com/main
user=ubuntu
environment=NODE_ENV="production"
